# deployment/systemd/nija-alerts.service
[Unit]
Description=NIJA Status Alerts Check
After=network.target

[Service]
Type=oneshot
User=nija
WorkingDirectory=/opt/nija
Environment="DATABASE_URL=postgresql://nija:password@localhost/nija"
ExecStart=/usr/bin/python3 /opt/nija/status_alerts.py
