# deployment/systemd/nija-alerts.timer
[Unit]
Description=NIJA Status Alerts Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
